# s6-overlay/s6-rc.d/code-server/run
#!/command/with-contenv bash
set -euo pipefail

exec s6-setuidgid agent code-server --bind-addr 0.0.0.0:8080 --auth none /home/agent/workspace
