# s6-overlay/s6-rc.d/system-bootstrap/up
/usr/local/share/copilot-code-server-container/system-bootstrap.sh
