# s6-overlay/s6-rc.d/agent-bootstrap/up
/command/with-contenv s6-setuidgid agent /usr/local/share/copilot-code-server-container/agent-bootstrap.sh true
